# .devcontainer/_build-context/rootfs/etc/supervisor/conf.d.available/desktop-lite.conf
[program:desktop-lite] ; 定义名为 desktop-lite 的受管服务
command=/usr/local/bin/h-service-desktop-lite ; 启动桌面服务的入口脚本
directory=/home/usr_vscode ; 进程工作目录
user=usr_vscode ; 以普通用户运行，避免长期以 root 身份跑桌面进程
autostart=true ; supervisord 启动后自动拉起该服务
autorestart=true ; 进程异常退出后自动重启
startsecs=5 ; 连续存活 5 秒后，才算启动成功
startretries=3 ; 启动失败时最多重试 3 次
stopsignal=TERM ; 停止服务时先发送 SIGTERM
stopasgroup=true ; 停止时把信号发给整个进程组，避免子进程残留
killasgroup=true ; 强制杀进程时也作用到整个进程组
redirect_stderr=true ; 把标准错误合并到标准输出日志
stdout_logfile=/var/log/supervisor/desktop-lite.log ; 标准输出日志文件
stderr_logfile=/var/log/supervisor/desktop-lite.err ; 预留的标准错误日志文件，当前 redirect_stderr=true 时通常不会单独使用
# stdout_logfile=/dev/fd/1 ; 如果要把标准输出直接打到容器 stdout，可改用这一行
# stderr_logfile=/dev/fd/2 ; 如果要把标准错误直接打到容器 stderr，可改用这一行
stdout_logfile_maxbytes=10MB ; 单个标准输出日志文件最大 10MB，设为 0 表示不限制
stdout_logfile_backups=3 ; 标准输出日志轮转时保留 3 个历史文件
environment=HOME="/home/usr_vscode",USER="usr_vscode",LOGNAME="usr_vscode" ; 显式修正用户态环境变量，避免 setuid 后仍沿用 root 的 USER/LOGNAME
