# .devcontainer/_build-context/rootfs/etc/supervisor/conf.d.available/opencode.conf
[program:opencode]
command=/usr/local/bin-priority/opencode serve --port 4096 --hostname 0.0.0.0 ; --print-logs --log-level=DEBUG
directory=/home/usr_vscode
user=usr_vscode
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/opencode.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/home/usr_vscode",USER="usr_vscode",LOGNAME="usr_vscode"
